# systemd/system/cloudflared-update.timer
[Unit]
Description=Update Cloudflare Tunnel

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
