# xdg/autostart/msm_notifier.desktop
[Desktop Entry]
Exec=msm_notifier
Name=MSM Notifier
Type=Application
StartupNotify=false
Terminal=false
Hidden=false
